# bin/minecraftserver40
#!/bin/bash
# /etc/init.d/minecraftserver40
# version 1.0.0 2015-07-04 (YYYY-MM-DD)

### BEGIN INIT INFO
# Provides:   minecraftserver40
# Required-Start: $local_fs $remote_fs $network $syslog $named
# Required-Stop:  $local_fs $remote_fs $network $syslog $named
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description:    Minecraft server 40
# Description:    Starts the minecraft server 40
### END INIT INFO

# Set up the common Settings
. /proj/minecraft/bin/settings40.sh

# Now that it is all set up do the real work
. /proj/minecraft/bin/minecraftserver.sh
